# cmpt.conf
start on startup
stop on shutdown

respawn
exec python /home/ubuntu/cmpt/worker.py

